# search.ini
[program:search]
directory=/var/search
command=/var/search/venv/bin/python main.py
autostart=true
autorestart=true
stdout_logfile=/var/log/search-stdout.log
stderr_logfile=/var/log/search-stderr.log
startretries=3
