# search.ini
[program:search]
directory=/var/search
command=/var/search/venv/bin/python main.py
autostart=true
autorestart=true
startretries=3
